# etc/inittab
::sysinit:/etc/rc.d/start
ttyS1::respawn:/bin/sh
